# vendor/plugins/factory_utils/script/solr
#! /bin/bash
# Comments to support chkconfig on RedHat Linux
# chkconfig: 2345 65 37
# description: Search engine server (a wrapper around Lucene)

# /etc/init.d/acts_as_solr: start and stop the solr server

# Source function library.      This creates the operating environment for the process to be started
. /etc/rc.d/init.d/functions
case "$1" in
 start)
       echo "Starting SOLR server..."
       cd /var/www/apps/{project_dir}/current/vendor/plugins/acts_as_solr/solr

       if [ "$USER" = "deploy" ] ; then
         exec /usr/bin/java -Dsolr.data.dir=/var/www/apps/{project_dir}/shared/solr/data/production -Djetty.port=8983 -jar start.jar &> /var/www/apps/mania_tv_admin/shared/log/solr.log &
       else
         exec su -c "/usr/bin/java -Dsolr.data.dir=/var/www/apps/{project_dir}/shared/solr/data/production -Djetty.port=8983 -jar start.jar &> /var/www/apps/{project_dir}/shared/log/solr.log &" deploy
       fi

       echo started
       ;;

 stop)
       echo "Stopping SOLR server..."
       ps ax | grep java | grep solr | grep -v grep | sed 's/^[ \t]*//' | cut -f 1 -d " " | xargs kill
       echo stopped
       ;;
 restart)
       $0 stop
       echo "Waiting for processes to stop..."
       sleep 3
       $0 start
       echo restarted
       ;;
 *)
       echo "Usage: /etc/init.d/solr {start|stop|restart}"
       exit 1
esac

exit 0